# flask-uwsgi.conf
start on [2345]
stop on [06]

script
    cd /path/to/flask-uwsgi
    exec uwsgi --ini /etc/flask-uwsgi/flask-uwsgi.ini
end script
